# pi/file_system/etc/xdg/autostart/scorpion.desktop
[Desktop Entry]
Name=scorpion
Exec=/home/ben/run.sh

